# deploy/loopctl.service
[Unit]
Description=loopctl - Agent-native project state store
Documentation=https://github.com/mkreyman/loopctl
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=loopctl
Group=loopctl
WorkingDirectory=/opt/loopctl
EnvironmentFile=/opt/loopctl/.env

ExecStart=/opt/loopctl/_build/prod/rel/loopctl/bin/loopctl start
ExecStop=/opt/loopctl/_build/prod/rel/loopctl/bin/loopctl stop

Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loopctl

[Install]
WantedBy=multi-user.target
